# mombot.service
[Unit]
Description=MomBot (Piglet) Voice Companion
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=will
Group=will
WorkingDirectory=/home/will/mom-bot
# Kill any stale process holding GPIO before we start
ExecStartPre=-/usr/bin/fuser -k /dev/gpiochip4
ExecStartPre=-/usr/bin/fuser -k /dev/gpiochip0
ExecStartPre=-/bin/sleep 1
ExecStart=/home/will/mom-bot/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

# Give the process time to clean up GPIO, then force kill
TimeoutStopSec=8
KillMode=mixed
KillSignal=SIGTERM

# GPIO + SPI + audio access
SupplementaryGroups=spi gpio audio video i2c

# Environment for ALSA and locale
Environment=HOME=/home/will
Environment=PYTHONUNBUFFERED=1
Environment=TERM=xterm

# Logging to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mombot

[Install]
WantedBy=multi-user.target
