# mom-bot/mombot.service
[Unit]
Description=MomBot (Piglet) Voice Companion
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/mom-bot/newApp
ExecStart=/home/pi/mom-bot/newApp/.venv/bin/python main.py

# On stop: send SIGTERM first (triggers atexit cleanup + GPIO release).
# If still alive after 10s, SIGKILL. Then wait 2s before any restart
# so the kernel fully releases GPIO handles.
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=3

# GPIO + SPI + audio access
SupplementaryGroups=spi gpio audio video i2c

# Environment for ALSA and locale
Environment=HOME=/home/pi
Environment=PYTHONUNBUFFERED=1
Environment=TERM=xterm

# Logging to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mombot

[Install]
WantedBy=multi-user.target
